# .config/autostart/localxrdb.desktop
[Desktop Entry] 
Encoding=UTF-8
Name=LocalXrdb
Comment=Load local xrdb
TryExec=localxrdb.sh
Exec=localxrdb.sh
