# .config/autostart/wpa_cli.desktop
[Desktop Entry] 
Encoding=UTF-8
Name=WpaCliPass
Comment=WpaCliPass
TryExec=wpaclipass.sh
Exec=wpaclipass.sh
